# packaging/com.glowbar.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent for the native glowbar Touch Bar app.

  Install:
    1. Build the app:   cd app && swift build -c release
    2. Copy this file, replacing __GLOWBAR_BIN__ with the built binary path
       (e.g. /Users/you/Projects/.../app/.build/release/glowbar-touchbar):
         sed "s#__GLOWBAR_BIN__#$PWD/app/.build/release/glowbar-touchbar#" \
           packaging/com.glowbar.app.plist > ~/Library/LaunchAgents/com.glowbar.app.plist
    3. launchctl load ~/Library/LaunchAgents/com.glowbar.app.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.glowbar.app.plist
    rm ~/Library/LaunchAgents/com.glowbar.app.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.glowbar.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>__GLOWBAR_BIN__</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>StandardErrorPath</key>
    <string>/tmp/glowbar.err.log</string>
    <key>StandardOutPath</key>
    <string>/tmp/glowbar.out.log</string>
</dict>
</plist>
